# nas/rclone-sharepoint.service
[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=/usr/bin/rclone serve restic -v sharepoint-backup:nas

[Install]
WantedBy=default.target
